# tomcat/tomcat-supervisor.conf
[supervisord]
nodaemon = true

[program: tomcat]
user = tomcat
command = /bin/bash -c "/usr/local/bin/run-tomcat.sh"