# kibana/kibana-supervisor.conf
[supervisord]
nodaemon = true

[program: kibana]
command = /bin/bash -c "/usr/local/bin/run-kibana.sh"
user = kibana

