# stacks/nodejs-app/app/app-supervisor.conf
[supervisord]
nodaemon = true

[program: http-server-8080]
command = /bin/bash -c "exec nodejs /usr/local/app/bin/http-server.js 8080"
stdout_logfile = /var/log/app/http-server-8080.log
