# apache2/apache2-supervisor.conf
[supervisord]
nodaemon = true

[program: apache2]
command = /bin/bash -c "/usr/local/bin/run-apache2.sh"