# mongodb/mongodb-supervisor.conf
[supervisord]
nodaemon = true

[program: mongodb]
command = /bin/bash -c "/usr/local/bin/run-mongodb.sh"