# elasticsearch/elasticsearch-supervisor.conf
[supervisord]
nodaemon = true

[program: elasticsearch]
command = /bin/bash -c "/usr/local/bin/run-elasticsearch.sh"
user = es

