# nginx/nginx-supervisor.conf
[supervisord]
nodaemon = true

[program: nginx]
command = /bin/bash -c "/usr/local/bin/run-nginx.sh"
