# couchdb/couchdb-supervisor.conf
[supervisord]
nodaemon = true

[program: couchdb]
command = /bin/bash -c "/usr/local/bin/run-couchdb.sh"
stdout_logfile = /var/log/supervisor/couchdb.log
user = couchdb

