# logstash/logstash-supervisor.conf
[supervisord]
nodaemon = true

[program: logstash]
command = /bin/bash -c "/usr/local/bin/run-elasticsearch.sh"
user = logstash

