# redis/redis-supervisor.conf
[supervisord]
nodaemon = true

[program: redis]
command = /bin/bash -c "/usr/local/bin/run-redis.sh"
stdout_logfile = /var/log/supervisor/redis.log
user = redis

